# roles/app-build/templates/student.service.j2
[Unit]
Description=Ansible Student Service

[Service]
ExecStart=gradle clean bootrun {{ TOMCAT_APP_HOME}}/student
StandardOutput=null

[Install]
WantedBy=multi-user.target
Alias=student.service
